# output/packs/incident-response.conf
{
  "queries": {
    "alf": {
      "description": "Retrieves the configuration values for the Application Layer Firewall for OSX.",
      "interval": "3600",
      "platform": "darwin",
      "query": "select * from alf;",
      "value": "Verify firewall settings are as restrictive as you need. Identify unwanted firewall holes made by malware or humans",
      "version": "1.4.5"
    },
    "alf_exceptions": {
      "description": "Retrieves the exceptions for the Application Layer Firewall in OSX.",
      "interval": "3600",
      "platform": "darwin",
      "query": "select * from alf_exceptions;",
      "value": "Verify firewall settings are as restrictive as you need. Identify unwanted firewall holes made by malware or humans",
      "version": "1.4.5"
    },
    "alf_explicit_auths": {
      "description": "Retrieves the list of processes with explicit authorization for the Application Layer Firewall.",
      "interval": "3600",
      "platform": "darwin",
      "query": "select * from alf_explicit_auths;",
      "value": "Verify firewall settings are as restrictive as you need. Identify unwanted firewall holes made by malware or humans",
      "version": "1.4.5"
    },
    "alf_services": {
      "description": "Retrieves the services for the Application Layer Firewall in OSX.",
      "interval": "3600",
      "platform": "darwin",
      "query": "select * from alf_services;",
      "value": "Verify firewall settings are as restrictive as you need. Identify unwanted firewall holes made by malware or humans",
      "version": "1.4.5"
    },
    "app_schemes": {
      "description": "Retrieves the list of application scheme/protocol-based IPC handlers.",
      "interval": "86400",
      "platform": "darwin",
      "query": "select * from app_schemes;",
      "value": "Post-priori hijack detection, detect potential sensitive information leakage.",
      "version": "1.4.7"
    },
    "arp_cache": {
      "description": "Retrieves the ARP cache values in the target system.",
      "interval": "3600",
      "query": "select * from arp_cache;",
      "value": "Determine if MITM in progress.",
      "version": "1.4.5"
    },
    "crontab": {
      "description": "Retrieves all the jobs scheduled in crontab in the target system.",
      "interval": "3600",
      "platform": "posix",
      "query": "select * from crontab;",
      "value": "Identify malware that uses this persistence mechanism to launch at a given interval",
      "version": "1.4.5"
    },
    "disk_encryption": {
      "description": "Retrieves the current disk encryption status for the target system.",
      "interval": "86400",
      "platform": "posix",
      "query": "select * from disk_encryption;",
      "value": "Identifies a system potentially vulnerable to disk cloning.",
      "version": "1.4.5"
    },
    "etc_hosts": {
      "description": "Retrieves all the entries in the target system /etc/hosts file.",
      "interval": "86400",
      "platform": "posix",
      "query": "select * from etc_hosts;",
      "value": "Identify network communications that are being redirected. Example: identify if security logging has been disabled",
      "version": "1.4.5"
    },
    "installed_applications": {
      "description": "Retrieves all the currently installed applications in the target OSX system.",
      "interval": "3600",
      "platform": "darwin",
      "query": "select * from apps;",
      "value": "Identify malware, adware, or vulnerable packages that are installed as an application.",
      "version": "1.4.5"
    },
    "ip_forwarding": {
      "description": "Retrieves the current status of IP/IPv6 forwarding.",
      "interval": "3600",
      "platform": "posix",
      "query": "select * from system_controls where oid = '4.30.41.1' union select * from system_controls where oid = '4.2.0.1';",
      "value": "Identify if a machine is being used as relay.",
      "version": "1.4.5"
    },
    "iptables": {
      "description": "Retrieves the current filters and chains per filter in the target system.",
      "interval": "3600",
      "platform": "linux",
      "query": "select * from iptables;",
      "value": "Verify firewall settings are as restrictive as you need. Identify unwanted firewall holes made by malware or humans",
      "version": "1.4.5"
    },
    "kernel_modules": {
      "description": "Retrieves all the information for the current kernel modules in the target Linux system.",
      "interval": "3600",
      "platform": "linux",
      "query": "select * from kernel_modules;",
      "value": "Identify malware that has a kernel module component.",
      "version": "1.4.5"
    },
    "kextstat": {
      "description": "Retrieves all the information about the current kernel extensions for the target OSX system.",
      "interval": "3600",
      "platform": "darwin",
      "query": "select * from kernel_extensions;",
      "value": "Identify malware that has a kernel extension component.",
      "version": "1.4.5"
    },
    "last": {
      "description": "Retrieves the list of the latest logins with PID, username and timestamp.",
      "interval": "3600",
      "platform": "posix",
      "query": "select * from last;",
      "value": "Useful for intrusion detection and incident response. Verify assumptions of what accounts should be accessing what systems and identify machines accessed during a compromise.",
      "version": "1.4.5"
    },
    "launchd": {
      "description": "Retrieves all the daemons that will run in the start of the target OSX system.",
      "interval": "3600",
      "platform": "darwin",
      "query": "select * from launchd;",
      "value": "Identify malware that uses this persistence mechanism to launch at system boot",
      "version": "1.4.5"
    },
    "listening_ports": {
      "description": "Retrieves all the listening ports in the target system.",
      "interval": "3600",
      "platform": "posix",
      "query": "select * from listening_ports;",
      "value": "Detect if a listening port iis not mapped to a known process. Find backdoors.",
      "version": "1.4.5"
    },
    "logged_in_users": {
      "description": "Retrieves the list of all the currently logged in users in the target system.",
      "interval": "3600",
      "platform": "posix",
      "query": "select liu.*, p.name, p.cmdline, p.cwd, p.root from logged_in_users liu, processes p where liu.pid = p.pid;",
      "value": "Useful for intrusion detection and incident response. Verify assumptions of what accounts should be accessing what systems and identify machines accessed during a compromise.",
      "version": "1.4.5"
    },
    "loginwindow1": {
      "description": "Retrieves all the values for the loginwindow process in the target OSX system.",
      "interval": "86400",
      "platform": "darwin",
      "query": "select key, subkey, value from plist where path = '/Library/Preferences/com.apple.loginwindow.plist';",
      "value": "Identify malware that uses this persistence mechanism to launch at system boot",
      "version": "1.4.5"
    },
    "loginwindow2": {
      "description": "Retrieves all the values for the loginwindow process in the target OSX system.",
      "interval": "86400",
      "platform": "darwin",
      "query": "select key, subkey, value from plist where path = '/Library/Preferences/loginwindow.plist';",
      "value": "Identify malware that uses this persistence mechanism to launch at system boot",
      "version": "1.4.5"
    },
    "loginwindow3": {
      "description": "Retrieves all the values for the loginwindow process in the target OSX system.",
      "interval": "86400",
      "platform": "darwin",
      "query": "select username, key, subkey, value from plist p, (select * from users where directory like '/Users/%') u where p.path = u.directory || '/Library/Preferences/com.apple.loginwindow.plist';",
      "value": "Identify malware that uses this persistence mechanism to launch at system boot",
      "version": "1.4.5"
    },
    "loginwindow4": {
      "description": "Retrieves all the values for the loginwindow process in the target OSX system.",
      "interval": "86400",
      "platform": "darwin",
      "query": "select username, key, subkey, value from plist p, (select * from users where directory like '/Users/%') u where p.path = u.directory || '/Library/Preferences/loginwindow.plist';",
      "value": "Identify malware that uses this persistence mechanism to launch at system boot",
      "version": "1.4.5"
    },
    "mounts": {
      "description": "Retrieves the current list of mounted drives in the target system.",
      "interval": "3600",
      "platform": "posix",
      "query": "select * from mounts;",
      "value": "Scope for lateral movement. Potential exfiltration locations. Potential dormant backdoors.",
      "version": "1.4.5"
    },
    "nfs_shares": {
      "description": "Retrieves the current list of Network File System mounted shares.",
      "interval": "3600",
      "platform": "darwin",
      "query": "select * from nfs_shares;",
      "value": "Scope for lateral movement. Potential exfiltration locations. Potential dormant backdoors.",
      "version": "1.4.5"
    },
    "open_files": {
      "description": "Retrieves all the open files per process in the target system.",
      "interval": "86400",
      "platform": "posix",
      "query": "select distinct pid, path from process_open_files where path not like '/private/var/folders%' and path not like '/System/Library/%' and path not in ('/dev/null', '/dev/urandom', '/dev/random');",
      "value": "Identify processes accessing sensitive files they shouldn't",
      "version": "1.4.5"
    },
    "open_sockets": {
      "description": "Retrieves all the open sockets per process in the target system.",
      "interval": "86400",
      "platform": "posix",
      "query": "select distinct pid, family, protocol, local_address, local_port, remote_address, remote_port, path from process_open_sockets where path <> '' or remote_address <> '';",
      "value": "Identify malware via connections to known bad IP addresses as well as odd local or remote port bindings",
      "version": "1.4.5"
    },
    "process_env": {
      "description": "Retrieves all the environment variables per process in the target system.",
      "interval": "86400",
      "platform": "posix",
      "query": "select * from process_envs;",
      "value": "Insight into the process data: Where was it started from, was it preloaded...",
      "version": "1.4.5"
    },
    "process_memory": {
      "description": "Retrieves the memory map per process in the target Linux system.",
      "interval": "86400",
      "platform": "linux",
      "query": "select * from process_memory_map;",
      "value": "Ability to compare with known good. Identify mapped regions corresponding with or containing injected code.",
      "version": "1.4.5"
    },
    "ramdisk": {
      "description": "Retrieves all the ramdisk currently mounted in the target system.",
      "interval": "3600",
      "platform": "posix",
      "query": "select * from block_devices where type = 'Virtual Interface';",
      "value": "Identify if an attacker is using temporary, memory storage to avoid touching disk for anti-forensics purposes",
      "version": "1.4.5"
    },
    "recent_items": {
      "description": "Retrieves the list of recent items opened in OSX by parsing the plist per user.",
      "interval": "86400",
      "platform": "darwin",
      "query": "select username, key, value from plist p, (select * from users where directory like '/Users/%') u where p.path = u.directory || '/Library/Preferences/com.apple.recentitems.plist';",
      "value": "Identify recently accessed items. Useful for compromised hosts.",
      "version": "1.4.5"
    },
    "sandboxes": {
      "description": "Lists the application bundle that owns a sandbox label.",
      "interval": "86400",
      "platform": "darwin",
      "query": "select * from sandboxes;",
      "value": "Post-priori hijack detection, detect potential sensitive information leakage.",
      "version": "1.4.7"
    },
    "shell_history": {
      "description": "Retrieves the command history, per user, by parsing the shell history files.",
      "interval": "86400",
      "platform": "posix",
      "query": "select * from users join shell_history using (uid);",
      "value": "Identify actions taken. Useful for compromised hosts.",
      "version": "1.4.5"
    },
    "startup_items": {
      "description": "Retrieve all the items that will load when the target OSX system starts.",
      "interval": "86400",
      "platform": "darwin",
      "query": "select * from startup_items;",
      "value": "Identify malware that uses this persistence mechanism to launch at a given interval",
      "version": "1.4.5"
    },
    "suid_bin": {
      "description": "Retrieves all the files in the target system that are setuid enabled.",
      "interval": "3600",
      "platform": "posix",
      "query": "select * from suid_bin;",
      "value": "Detect backdoor binaries (attacker may drop a copy of /bin/sh). Find potential elevation points / vulnerabilities in the standard build.",
      "version": "1.4.5"
    },
    "wireless_networks": {
      "description": "Retrieves all the remembered wireless network that the target machine has connected to.",
      "interval": "3600",
      "platform": "darwin",
      "query": "select ssid, network_name, security_type, last_connected, captive_portal, possibly_hidden, roaming, roaming_profile from wifi_networks;",
      "value": "Identifies connections to rogue access points.",
      "version": "1.6.0"
    }
  }
}
